# services/inexcons-backend.service
[Unit]
Description=INEXCONS Laravel Backend Service
After=network.target mysql.service
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/mnt/volume_nyc1_01/inexcons/backend
ExecStart=/usr/bin/php artisan serve --host=0.0.0.0 --port=8000
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=5
TimeoutStopSec=30

# Environment variables
Environment=PHP_CLI_SERVER_WORKERS=4
Environment=APP_ENV=production

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/mnt/volume_nyc1_01/inexcons/backend/storage
ReadWritePaths=/mnt/volume_nyc1_01/inexcons/backend/bootstrap/cache
ReadWritePaths=/mnt/volume_nyc1_01/inexcons/backend/database

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=inexcons-backend

[Install]
WantedBy=multi-user.target
